# debian/kodi-rockchip-gbm.postinst
#!/bin/sh
set -e

#DEBHELPER#

echo "kodi-rockchip-gbm: installed. Run 'sudo systemctl start kodi' to run, and 'sudo systemctl enable kodi' to enable at boot."

exit 0
